# examples/saving_photos_for_relearning/supervisord.conf
[supervisord]
nodaemon = false

[inet_http_server]
port = 9001

[supervisorctl]

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:source]
command = "python.exe" "%(here)s/../../owl/module_source_cv.py" "%(here)s/config.json"
process_name = cv_%(process_num)s
numprocs = 1
stopsignal=TERM
autostart=true
autorestart=true

[program:sifter]
command = "python.exe" "%(here)s/../../owl/module_sifter.py" "%(here)s/config.json"
process_name = sifter_%(process_num)s
numprocs = 1
stopsignal=TERM
autostart=true
autorestart=true

[program:blur_check]
command = "python.exe" "%(here)s/../../owl/module_blur_check.py" "%(here)s/config.json"
process_name = blur_check_%(process_num)s
numprocs = 1
stopsignal=TERM
autostart=true
autorestart=true

[program:save_photos]
command = "python.exe" "%(here)s/../../owl/module_save_photos.py" "%(here)s/config.json"
process_name = save_photos_%(process_num)s
numprocs = 1
stopsignal=TERM
autostart=true
autorestart=true

[program:sink]
command = "python.exe" "%(here)s/../../owl/module_sink_win.py" "%(here)s/config.json"
process_name = win_%(process_num)s
numprocs = 1
stopsignal=TERM
autostart=true
autorestart=true
