# examples/file_file/supervisord.conf
[supervisord]
nodaemon = false

[inet_http_server]
port = 9001

[supervisorctl]

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:source]
command = "c:/Python38/python.exe" "%(here)s/../../owl/module_source_dshow.py" "%(here)s/config.json"
process_name = dshow_%(process_num)s
numprocs = 1
stopsignal=TERM
autostart=true
autorestart=true

[program:sink]
command = "c:/Python38/python.exe" "%(here)s/../../owl/module_sink_file.py" "%(here)s/config.json"
process_name = sink_file_%(process_num)s
numprocs = 3
stopsignal=TERM
autostart=true
autorestart=true