# examples/slicer/supervisordl.conf
[supervisord]
nodaemon = false

[inet_http_server]
port = 9001

[supervisorctl]

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:source]
command = "python3" "%(here)s/../../owl/module_source_cv.py" "%(here)s/config.json"
process_name = cv_%(process_num)s
numprocs = 1
stopsignal=TERM
autostart=true
autorestart=true

[program:process_slice]
command = "python3" "%(here)s/../../owl/module_slice.py" "%(here)s/config.json"
process_name = slice_%(process_num)s
numprocs = 1
stopsignal=TERM
autostart=true
autorestart=true

[program:sink]
command = "python3" "%(here)s/../../owl/module_sink_win.py" "%(here)s/config.json"
process_name = win_%(process_num)s
numprocs = 2
stopsignal=TERM
autostart=true
autorestart=true